# processing-service/systemd/meeting-processor.user.service
[Unit]
Description=Private Meeting Processing Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kubuntiai/meeting-transcriber-processing
EnvironmentFile=/home/kubuntiai/meeting-transcriber-processing/.env
Environment=HF_HUB_OFFLINE=1
Environment=HF_HUB_DISABLE_TELEMETRY=1
Environment=TRANSFORMERS_OFFLINE=1
ExecStart=/home/kubuntiai/meeting-transcriber-processing/.venv/bin/uvicorn meeting_processor.app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=default.target
